# Resources/Contents/Library/LaunchDaemons/com.batterycare.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.batterycare.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/BatteryCare.app/Contents/MacOS/battery-care-daemon</string>
    </array>

    <key>UserName</key>
    <string>root</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Library/Logs/BatteryCare/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Logs/BatteryCare/daemon-error.log</string>
</dict>
</plist>
